# deployment/jlpt-celery-worker.service
[Unit]
Description=JLPT Celery Worker
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/mikan
Environment="PATH=/var/www/mikan/venv/bin"
EnvironmentFile=/var/www/mikan/.env
Environment=DJANGO_SETTINGS_MODULE=config.settings.production
ExecStart=/var/www/mikan/venv/bin/celery -A config worker \
    --loglevel=info
ExecStop=/bin/kill -TERM $MAINPID
KillSignal=SIGTERM
KillMode=mixed
Restart=on-failure
RestartSec=10
# Log to journal (default). For file logs, use:
# ExecStart=... --logfile=/var/www/mikan/logs/celery_worker.log
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jlpt-celery-worker

[Install]
WantedBy=multi-user.target
